# platform-linux/replacement/.config/systemd/user/switch-external-display.timer
[Unit]
Description=Switch monitor each 1s

[Timer]
OnCalendar=*:*:0/1

[Install]
WantedBy=timers.target
