# platform-darwin/replacement/Library/LaunchAgents/localhost.dotfiles.launch_toggl_desktop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>localhost.dotfiles.launch_toggl_desktop</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>open</string>
    <string>-a</string>
    <string>TogglDesktop</string>
    <!-- I want to boot it at background -->
    <string>-gj</string>
  </array>
  <key>LaunchOnlyOnce</key>
  <true/>
</dict>
</plist>
